# deploy/hl7receiver.conf
description "hl7receiver"

start on runlevel [2345]
stop on runlevel [!2345]

setuid root
setgid root

env JAVA_HOME=/usr/lib/jvm/java-8-oracle/jre
env CONFIG_JDBC_USERNAME=endeavour
env CONFIG_JDBC_PASSWORD=
env CONFIG_JDBC_URL=jdbc:postgresql://127.0.0.1/config

exec java -jar /opt/hl7receiver/hl7receiver-1.0-SNAPSHOT-jar-with-dependencies.jar
